# Abbot.service
[Unit]
Description=Abbot Python Backend

[Service]
Type=simple
ExecStart=/home/pi/Desktop/Abbot/abbot_start.py

[Install]
WantedBy=multi-user.target
